# templates/default/butterfly.socket.erb
[Socket]
ListenStream=<%= @config['host'].gsub(/^\"|\"?$/, '') %>:<%= @config['port'] %>

[Install]
WantedBy=sockets.target
