# templates/default/butterfly.socket.erb
[Socket]
ListenStream=<%= node['butterfly']['config']['host'].rchomp('"').chomp('"') %>:<%= node['butterfly']['config']['port'] %>

[Install]
WantedBy=sockets.target
